# supervisor/carbon.conf
[program:carbon]
command=carbon-cache.py --nodaemon start
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
